# data/packagekit.service.in
[Unit]
Description=PackageKit Daemon
# PK doesn't know how to do anything on ostree-managed systems;
# currently the design is to have dedicated daemons like
# eos-updater and rpm-ostree, and gnome-software talks to those.
ConditionPathExists=!/run/ostree-booted
# PK shouldn't run while librpm is updated
ConditionPathExists=!/var/lib/PackageKit/upgrade_lock

[Service]
Type=dbus
BusName=org.freedesktop.PackageKit
User=@PACKAGEKIT_USER@
ExecStart=@libexecdir@/packagekitd
# don't kill service immediately,
# instead request it to quit after finishing current request
ExecStop=@bindir@/pkcon quit
